# supervisord.conf
[supervisord]
nodaemon=true

[program:apache2]
command=/usr/sbin/apache2ctl -D FOREGROUND

[program:mailcatcher]
command=/usr/local/bin/mailcatcher --foreground --ip=0.0.0.0